# system-scripts/etc/cron.weekly/btrfs-clean-job
#!/usr/bin/env bash

set -e

btrfs scrub start /
