# system-scripts/etc/cron.weekly/clean-pkg-managers
#!/bin/bash

# Flatpak Cleanup and Repair
echo "--- Running Flatpak Cleanup and Repair ---"
flatpak uninstall --unused -y
flatpak repair

# Pacman Orphan Removal
echo "--- Running Pacman Orphan Removal ---"
ORPHANS=$(pacman -Qdtq)
if [ -n "$ORPHANS" ]; then
    sudo pacman -Rns --noconfirm "$ORPHANS"
else
    echo "No orphaned packages found."
fi

# Pacman Cache Clear
echo "--- Running Pacman Cache Prune (Keeping last 3) ---"
sudo paccache -r
